# files/spark-uspto
#!/bin/sh
#
# /etc/init.d/spark
# init script for Spark.
# generated by 'spark enable boot-start'.
#
# chkconfig: 2345 90 60
# description: Spark indexer service
#
RETVAL=0


spark_start() {
  echo Starting Spark...
  "/opt/Spark/spark-2.1.1-bin-hadoop2.7/sbin/start-all.sh" start --no-prompt --answer-yes
}
spark_stop() {
  echo Stopping Spark...
  "/opt/Spark/spark-2.1.1-bin-hadoop2.7/sbin/start-all.sh" stop
}
spark_restart() {
  echo Restarting Spark...
  "/opt/Spark/spark-2.1.1-bin-hadoop2.7/sbin/start-all.sh" restart
}
spark_status() {
  echo Spark status:
  "/opt/Spark/spark-2.1.1-bin-hadoop2.7/sbin/start-all.sh" status
}
case "$1" in
  start)
    spark_start
    ;;
  stop)
    spark_stop
    ;;
  restart)
    spark_restart
    ;;
  status)
    spark_status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit $RETVAL

